# systemd-service/autoscreen.timer
[Unit]
Description=Run autoscreen every 15 minutes.

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
