# systemd-service/autoscreen.timer
[Unit]
Description=Run autoscreen every hour.

[Timer]
OnCalendar=hourly
RandomizedDelaySec=3600
AccuracySec=1us

[Install]
WantedBy=timers.target
