# systemd-service/autoscreen.timer
[Unit]
Description=Run autoscreen every 15 minutes.

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
